# var/lib/opkg/info/onecloudupgrade.ipk.list
/usr/share/onecloud-upgrade/upgrade_rootfs.sh
/usr/share/onecloud-upgrade/opkg.conf
/usr/share/onecloud-upgrade/upgrade_log_debug.conf
/usr/share/onecloud-upgrade/upgrade_log.conf
/usr/share/onecloud-upgrade/.opkg_default.conf
/usr/share/onecloud-upgrade/.img_date
/usr/sbin/onecloud_upgrading_led.sh
/usr/sbin/onecloud_upgrade_run.sh
/usr/sbin/opkg-extract-file
/usr/sbin/upgrade
/etc/init.d/S88wxmnt
